# zoom-update.cron
#!/bin/bash

# Resolve the full path of the script
SCRIPT_DIR=$(dirname "$(readlink -f "$0")")

exec "${SCRIPT_DIR}/zoom-update" "--cron"